# infra/bare-metal/crontech-deploy-agent.service
# Crontech Deploy Agent — internal privileged HTTP service
#
# Listens on http://127.0.0.1:9091 ONLY. Never exposed externally.
# Crontech's API admin endpoint (POST /api/admin/deploy) proxies here.
#
# Runs as root so it can: git pull, bun install, bun build, systemctl restart.
# The DEPLOY_AGENT_SECRET must match the value in /opt/crontech/.env.
#
# Install:
#   cp infra/bare-metal/crontech-deploy-agent.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now crontech-deploy-agent
#
# Hardening (Wave B2, 2026-04-27): added explicit PATH, OnFailure notifier,
# tightened restart limits, and Bun-friendly stop semantics. To enable
# systemd watchdog later, emit sd_notify WATCHDOG pings and uncomment
# WatchdogSec=60.

[Unit]
Description=Crontech Deploy Agent (internal, localhost:9091)
Documentation=https://api.crontech.ai
After=network.target
OnFailure=crontech-failure-notify@%i.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/crontech

ExecStart=/root/.bun/bin/bun run services/deploy-agent/src/index.ts
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10
KillSignal=SIGTERM
TimeoutStopSec=30s
KillMode=mixed
# WatchdogSec=60s  # enable once the Bun runtime emits sd_notify WATCHDOG pings

Environment=PATH=/usr/local/bin:/usr/bin:/bin:/root/.bun/bin
Environment=APP_DIR=/opt/crontech
Environment=DEPLOY_AGENT_PORT=9091
EnvironmentFile=-/opt/crontech/.env

# Only listen on localhost — belt-and-braces
IPAddressAllow=127.0.0.1/8 ::1/128
IPAddressDeny=any

StandardOutput=journal
StandardError=journal
SyslogIdentifier=crontech-deploy-agent

[Install]
WantedBy=multi-user.target
