# infra/bare-metal/crontech-edge-runtime.service
# Crontech Edge Runtime (BLK-017 v0)
#
# Self-hosted V8-isolate edge runtime. Bun.serve on 127.0.0.1:9096 ONLY.
# Internal — Caddy / tunnel daemon proxy customer traffic in.
#
# Per CLAUDE.md §0.10 Zero-Idle Rule: service unit lives in repo so
# deploy.yml auto-installs on every push to Main.
#
# Install (manual, if needed):
#   cp infra/bare-metal/crontech-edge-runtime.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now crontech-edge-runtime
#
# Hardening (Wave B2, 2026-04-27): added explicit PATH, OnFailure notifier,
# tightened restart limits, and Bun-friendly stop semantics. To enable
# systemd watchdog later, emit sd_notify WATCHDOG pings and uncomment
# WatchdogSec=60.

[Unit]
Description=Crontech Edge Runtime (V8-isolate, BLK-017 v0)
Documentation=https://docs.crontech.ai/blocks/blk-017
After=network.target
OnFailure=crontech-failure-notify@%i.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/crontech

ExecStart=/root/.bun/bin/bun run services/edge-runtime/src/index.ts
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10
KillSignal=SIGTERM
TimeoutStopSec=30s
KillMode=mixed
# WatchdogSec=60s  # enable once the Bun runtime emits sd_notify WATCHDOG pings

Environment=PATH=/usr/local/bin:/usr/bin:/bin:/root/.bun/bin
Environment=APP_DIR=/opt/crontech
Environment=EDGE_RUNTIME_PORT=9096
EnvironmentFile=-/opt/crontech/.env

# Only listen on localhost — belt-and-braces. Edge ingress comes through
# Caddy (TLS termination) or the tunnel daemon (BLK-019), never direct.
IPAddressAllow=127.0.0.1/8 ::1/128
IPAddressDeny=any

StandardOutput=journal
StandardError=journal
SyslogIdentifier=crontech-edge-runtime

[Install]
WantedBy=multi-user.target
